# etc/systemd/system/set-ryzenadj-tweaks.path
[Path]
PathChanged=/home/deck/.local/bin/allowadj.txt

[Install]
WantedBy=multi-user.target
